# wazuh/rules/custom_mitre_mapped.xml
<!-- =============================================================================
  custom_mitre_mapped.xml — Regole con mapping MITRE ATT&CK esplicito
  ID range: 100060 – 100099
  Autore: SOC Home Lab

  Questo file raccoglie regole che wrappano/estendono detection built-in di Wazuh
  aggiungendo tag <mitre> strutturati per le tecniche principali del lab.
  L'obiettivo è avere un mapping pulito e navigabile dalla Dashboard.

  Tecniche coperte:
    T1110   — Brute Force
    T1059   — Command and Scripting Interpreter
    T1053   — Scheduled Task/Job
    T1547   — Boot/Logon Autostart Execution
    T1021   — Remote Services
    T1562   — Impair Defenses
============================================================================= -->

<group name="custom_mitre,">

  <!-- =========================================================
    T1110 — BRUTE FORCE
    Aggrega: SSH brute force (Linux) + Windows login brute force
  ========================================================= -->

  <!-- Brute force SSH — estende la detection built-in di Wazuh -->
  <rule id="100060" level="10" frequency="6" timeframe="60">
    <if_matched_sid>5760</if_matched_sid>
    <same_source_ip />
    <description>[T1110.001] SSH Brute Force: multiple failures from $(srcip)</description>
    <mitre>
      <id>T1110</id>
      <id>T1110.001</id>
    </mitre>
    <group>brute_force,credential_access,pci_dss_10.2.4,</group>
  </rule>

  <!-- Brute force Windows — troppi logon failure (Event ID 4625) -->
  <rule id="100061" level="10" frequency="6" timeframe="60">
    <if_matched_sid>60122</if_matched_sid>
    <same_source_ip />
    <description>[T1110.001] Windows Brute Force: multiple logon failures from $(srcip)</description>
    <mitre>
      <id>T1110</id>
      <id>T1110.001</id>
    </mitre>
    <group>brute_force,credential_access,authentication_failures,</group>
  </rule>

  <!-- Password spray: stessa password, utenti diversi -->
  <rule id="100062" level="10" frequency="5" timeframe="60">
    <if_matched_sid>60122</if_matched_sid>
    <same_source_ip />
    <different_user />
    <description>[T1110.003] Password Spray: logon failures for different users from $(srcip)</description>
    <mitre>
      <id>T1110</id>
      <id>T1110.003</id>
    </mitre>
    <group>brute_force,credential_access,</group>
  </rule>


  <!-- =========================================================
    T1059 — COMMAND AND SCRIPTING INTERPRETER
  ========================================================= -->

  <!-- PowerShell offensivo (estende 100021 con MITRE tag esplicito) -->
  <rule id="100063" level="12">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)powershell\.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(-enc|-nop|-windowstyle\s+hidden|iex|invoke-expression|downloadstring)</field>
    <description>[T1059.001] PowerShell with offensive flags: $(win.eventdata.commandLine)</description>
    <mitre>
      <id>T1059</id>
      <id>T1059.001</id>
    </mitre>
    <group>execution,defense_evasion,</group>
  </rule>

  <!-- cmd.exe lanciato da processo non interattivo -->
  <rule id="100064" level="8">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)cmd\.exe</field>
    <field name="win.eventdata.parentImage" type="pcre2">(?i)(services|svchost|lsass|winlogon|csrss|smss)\.exe</field>
    <description>[T1059.003] cmd.exe spawned from system process $(win.eventdata.parentImage)</description>
    <mitre>
      <id>T1059</id>
      <id>T1059.003</id>
    </mitre>
    <group>execution,</group>
  </rule>

  <!-- Script bash/sh sospetto su Linux (via Wazuh audit) -->
  <rule id="100065" level="8">
    <if_sid>80784</if_sid>
    <field name="audit.exe" type="pcre2">(?i)(python|perl|ruby|nc|ncat|bash|sh)</field>
    <description>[T1059] Scripting interpreter executed: $(audit.exe) by $(audit.auid)</description>
    <mitre>
      <id>T1059</id>
    </mitre>
    <group>execution,</group>
  </rule>


  <!-- =========================================================
    T1053 — SCHEDULED TASK / JOB
  ========================================================= -->

  <!-- schtasks /create da linea di comando (Windows) -->
  <rule id="100070" level="10">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)schtasks\.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)/create</field>
    <description>[T1053.005] Scheduled task created via command line: $(win.eventdata.commandLine)</description>
    <mitre>
      <id>T1053</id>
      <id>T1053.005</id>
    </mitre>
    <group>persistence,execution,</group>
  </rule>

  <!-- Task schedulato su Linux via crontab -->
  <rule id="100071" level="8">
    <if_sid>5902</if_sid>
    <description>[T1053.003] Cron job modified by $(dstuser)</description>
    <mitre>
      <id>T1053</id>
      <id>T1053.003</id>
    </mitre>
    <group>persistence,</group>
  </rule>


  <!-- =========================================================
    T1547 — BOOT/LOGON AUTOSTART EXECUTION
  ========================================================= -->

  <!-- Registro Run key modificato -->
  <rule id="100080" level="12">
    <if_group>sysmon_event_13</if_group>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\(CurrentVersion\\Run|CurrentVersion\\RunOnce)</field>
    <description>[T1547.001] Registry Run key modified: $(win.eventdata.targetObject)</description>
    <mitre>
      <id>T1547</id>
      <id>T1547.001</id>
    </mitre>
    <group>persistence,</group>
  </rule>

  <!-- File creato in cartella Startup -->
  <rule id="100081" level="10">
    <if_group>sysmon_event_11</if_group>
    <field name="win.eventdata.targetFilename" type="pcre2">(?i)\\Programs\\Startup\\</field>
    <description>[T1547.001] File dropped in Startup folder: $(win.eventdata.targetFilename)</description>
    <mitre>
      <id>T1547</id>
      <id>T1547.001</id>
    </mitre>
    <group>persistence,</group>
  </rule>


  <!-- =========================================================
    T1021 — REMOTE SERVICES
  ========================================================= -->

  <!-- Login SSH riuscito (lateral movement via SSH) -->
  <rule id="100090" level="5">
    <if_sid>5715</if_sid>
    <description>[T1021.004] Successful SSH login: $(dstuser) from $(srcip)</description>
    <mitre>
      <id>T1021</id>
      <id>T1021.004</id>
    </mitre>
    <group>authentication_success,lateral_movement,</group>
  </rule>

  <!-- RDP connection (Windows Event ID 4624, logon type 10) -->
  <rule id="100091" level="6">
    <if_sid>60106</if_sid>
    <field name="win.eventdata.logonType" type="pcre2">^10$</field>
    <description>[T1021.001] Remote Desktop (RDP) login: $(win.eventdata.targetUserName) from $(win.eventdata.ipAddress)</description>
    <mitre>
      <id>T1021</id>
      <id>T1021.001</id>
    </mitre>
    <group>authentication_success,lateral_movement,</group>
  </rule>


  <!-- =========================================================
    T1562 — IMPAIR DEFENSES
  ========================================================= -->

  <!-- Disabilitazione Windows Defender via PowerShell -->
  <rule id="100095" level="14">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(Set-MpPreference|DisableRealtimeMonitoring|DisableAntiSpyware|sc\s+stop\s+windefend)</field>
    <description>[T1562.001] Windows Defender disabled or tampered: $(win.eventdata.commandLine)</description>
    <mitre>
      <id>T1562</id>
      <id>T1562.001</id>
    </mitre>
    <group>defense_evasion,</group>
  </rule>

  <!-- Cancellazione log eventi Windows -->
  <rule id="100096" level="14">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(Clear-EventLog|wevtutil\s+cl|wevtutil\s+clear)</field>
    <description>[T1562.002] Windows Event Log cleared: $(win.eventdata.commandLine)</description>
    <mitre>
      <id>T1562</id>
      <id>T1562.002</id>
    </mitre>
    <group>defense_evasion,</group>
  </rule>

</group>
