# wazuh/rules/custom_windows.xml
<!-- =============================================================================
  custom_windows.xml — Regole custom per detection Windows + Sysmon
  ID range: 100020 – 100059
  Autore: SOC Home Lab
  MITRE coverage: T1059 (Execution), T1547 (Persistence), T1055 (Injection),
                  T1204 (User Execution), T1003 (Credential Access)

  Prerequisiti:
    - Sysmon installato su target-windows con config SwiftOnSecurity
    - Wazuh Agent configurato per leggere Microsoft-Windows-Sysmon/Operational
    - Vedi docs/sysmon-setup.md

  Nota sui field name:
    I decoder Wazuh per eventi Windows espongono i campi con prefisso
    data.win.eventdata.* — nelle regole XML si usa <field> senza prefisso:
      <field name="win.eventdata.parentImage">...</field>
============================================================================= -->

<group name="custom_windows,sysmon,">

  <!-- =========================================================
    REGOLA 100020 — Process injection sospetto via Office
    Trigger: cmd.exe o powershell.exe spawned da Office (Sysmon ID 1)
    Tecnica: T1059.001 / T1566.001 (spearphishing attachment)
    Severity: 14 (Critical) — quasi sempre malevolo
  ========================================================= -->
  <rule id="100020" level="14">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.parentImage" type="pcre2">(?i)(winword|excel|powerpnt|outlook|mspub|msaccess)\.exe</field>
    <field name="win.eventdata.image" type="pcre2">(?i)(cmd|powershell|wscript|cscript|mshta|rundll32|regsvr32)\.exe</field>
    <description>Suspicious child process $(win.eventdata.image) spawned from Office application $(win.eventdata.parentImage)</description>
    <mitre>
      <id>T1059</id>
      <id>T1059.001</id>
      <id>T1566.001</id>
    </mitre>
    <group>execution,initial_access,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100021 — PowerShell con argomenti offensivi
    Trigger: powershell.exe con flag tipici di esecuzione offensiva
    Tecnica: T1059.001 (PowerShell)
    Severity: 12 (High)
  ========================================================= -->
  <rule id="100021" level="12">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)powershell\.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(-enc|-encodedcommand|-nop|-noprofile|-windowstyle\s+hidden|-exec\s+bypass|iex|invoke-expression|downloadstring|webclient)</field>
    <description>PowerShell execution with suspicious arguments: $(win.eventdata.commandLine)</description>
    <mitre>
      <id>T1059.001</id>
      <id>T1027</id>
    </mitre>
    <group>execution,defense_evasion,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100022 — Modifica chiavi Run/RunOnce nel registro
    Trigger: Sysmon Event ID 13 (Registry Value Set) su chiavi Run
    Tecnica: T1547.001 (Registry Run Keys / Startup Folder)
    Severity: 12 (High) — persistenza tramite autorun
  ========================================================= -->
  <rule id="100022" level="12">
    <if_group>sysmon_event_13</if_group>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\(CurrentVersion\\Run|CurrentVersion\\RunOnce|CurrentVersion\\RunOnceEx)</field>
    <description>Registry Run key modified for persistence: $(win.eventdata.targetObject) by $(win.eventdata.image)</description>
    <mitre>
      <id>T1547</id>
      <id>T1547.001</id>
    </mitre>
    <group>persistence,defense_evasion,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100023 — Creazione file eseguibile in path sospette
    Trigger: Sysmon Event ID 11 (File Create) in %TEMP% o %APPDATA%
    Tecnica: T1204.002 (Malicious File)
    Severity: 10 (High)
  ========================================================= -->
  <rule id="100023" level="10">
    <if_group>sysmon_event_11</if_group>
    <field name="win.eventdata.targetFilename" type="pcre2">(?i)\\(Temp|AppData\\Roaming|AppData\\Local\\Temp|ProgramData)\\.*\.(exe|dll|bat|ps1|vbs|js|hta|cmd)</field>
    <description>Executable file created in suspicious path: $(win.eventdata.targetFilename)</description>
    <mitre>
      <id>T1204.002</id>
      <id>T1059</id>
    </mitre>
    <group>execution,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100024 — Connessione di rete da processo non standard
    Trigger: Sysmon Event ID 3 (Network Connection) da processi
             insoliti che normalmente non fanno rete
    Tecnica: T1071 (Application Layer Protocol) / C2
    Severity: 10 (High)
  ========================================================= -->
  <rule id="100024" level="10">
    <if_group>sysmon_event3</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)(mshta|wscript|cscript|regsvr32|regsvcs|regasm|installutil|certutil|bitsadmin|forfiles|pcalua)\.exe</field>
    <description>Suspicious network connection initiated by $(win.eventdata.image) to $(win.eventdata.destinationIp):$(win.eventdata.destinationPort)</description>
    <mitre>
      <id>T1071</id>
      <id>T1218</id>
    </mitre>
    <group>command_and_control,defense_evasion,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100025 — LSASS memory access (credential dumping)
    Trigger: Sysmon Event ID 10 (Process Access) su lsass.exe
    Tecnica: T1003.001 (LSASS Memory)
    Severity: 14 (Critical) — quasi sempre credential dumping
  ========================================================= -->
  <rule id="100025" level="14">
    <if_group>sysmon_event_10</if_group>
    <field name="win.eventdata.targetImage" type="pcre2">(?i)lsass\.exe</field>
    <description>Potential credential dumping: process $(win.eventdata.sourceImage) accessed LSASS memory</description>
    <mitre>
      <id>T1003</id>
      <id>T1003.001</id>
    </mitre>
    <group>credential_access,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100026 — Scheduled Task creata da PowerShell o cmd
    Trigger: Sysmon Event ID 1, schtasks.exe con /create
    Tecnica: T1053.005 (Scheduled Task)
    Severity: 10 (High)
  ========================================================= -->
  <rule id="100026" level="10">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)schtasks\.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)/create</field>
    <description>Scheduled task created: $(win.eventdata.commandLine)</description>
    <mitre>
      <id>T1053</id>
      <id>T1053.005</id>
    </mitre>
    <group>persistence,execution,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100027 — DNS query verso dominio con DGA-like pattern
    Trigger: Sysmon Event ID 22 (DNS Query) — lunghezza nome elevata
    Tecnica: T1568 (Dynamic Resolution) — possibile C2 via DGA
    Severity: 8 (Medium) — richiede analisi contestuale
  ========================================================= -->
  <rule id="100027" level="8">
    <if_group>sysmon_event_22</if_group>
    <field name="win.eventdata.queryName" type="pcre2">(?i)[a-z0-9]{20,}\.(com|net|org|xyz|top|cc|pw)</field>
    <description>DNS query to potentially DGA-generated domain: $(win.eventdata.queryName) from $(win.eventdata.image)</description>
    <mitre>
      <id>T1568</id>
      <id>T1071.004</id>
    </mitre>
    <group>command_and_control,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100028 — Modifica script di startup Windows
    Trigger: Sysmon Event ID 11 — file creato in Startup folders
    Tecnica: T1547.001
    Severity: 10 (High)
  ========================================================= -->
  <rule id="100028" level="10">
    <if_group>sysmon_event_11</if_group>
    <field name="win.eventdata.targetFilename" type="pcre2">(?i)\\(Start Menu\\Programs\\Startup|ProgramData\\Microsoft\\Windows\\Start Menu\\Programs\\Startup)\\</field>
    <description>File created in Startup folder for persistence: $(win.eventdata.targetFilename)</description>
    <mitre>
      <id>T1547</id>
      <id>T1547.001</id>
    </mitre>
    <group>persistence,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100030 — Cancellazione Shadow Copies / backup (ransomware)
    Trigger: Sysmon EID 1, vssadmin/wmic/wbadmin/diskshadow che cancellano copie shadow
    Tecnica: T1490 (Inhibit System Recovery)
    Severity: 14 (Critical) — comportamento classico pre-cifratura ransomware
  ========================================================= -->
  <rule id="100030" level="14">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)(vssadmin|wmic|wbadmin|diskshadow)\.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(delete\s+shadows|shadowcopy\s+delete|delete\s+catalog|delete\s+systemstatebackup|resize\s+shadowstorage)</field>
    <description>[T1490] Shadow copy/backup deletion (ransomware pre-encryption): $(win.eventdata.commandLine)</description>
    <mitre>
      <id>T1490</id>
    </mitre>
    <group>impact,ransomware,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100031 — Disabilitazione recovery di Windows (bcdedit)
    Trigger: Sysmon EID 1, bcdedit che disabilita il recovery o ignora i failure
    Tecnica: T1490 (Inhibit System Recovery)
    Severity: 14 (Critical)
  ========================================================= -->
  <rule id="100031" level="14">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)bcdedit\.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(recoveryenabled\s+no|bootstatuspolicy\s+ignoreallfailures)</field>
    <description>[T1490] Windows boot recovery disabled via bcdedit (ransomware): $(win.eventdata.commandLine)</description>
    <mitre>
      <id>T1490</id>
    </mitre>
    <group>impact,ransomware,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100032 — Ransom note / file cifrati (estensioni note)
    Trigger: Sysmon EID 11, file creato con nome da ransom note o estensione di cifratura
    Tecnica: T1486 (Data Encrypted for Impact)
    Severity: 12 (High)
  ========================================================= -->
  <rule id="100032" level="12">
    <if_group>sysmon_event_11</if_group>
    <field name="win.eventdata.targetFilename" type="pcre2">(?i)(\\(readme|how[_\s-]?to[_\s-]?decrypt|decrypt[_\s-]?instruction|recover[_\s-]?(your[_\s-]?)?files|_?restore[_\s-]?files|ransom)[^\\]*\.(txt|html|hta)$|\.(locked|encrypted|crypt|crypted|enc|wcry|wncry|ryk|lockbit)$)</field>
    <description>[T1486] Possible ransomware artifact (note or encrypted file): $(win.eventdata.targetFilename)</description>
    <mitre>
      <id>T1486</id>
    </mitre>
    <group>impact,ransomware,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100033 — Dump della memoria LSASS via comsvcs/procdump/rundll32
    Trigger: Sysmon EID 1, command line tipica del dumping di LSASS
    Tecnica: T1003.001 (OS Credential Dumping: LSASS Memory)
    Severity: 14 (Critical) — complementa 100025 (che usa Sysmon EID 10)
  ========================================================= -->
  <rule id="100033" level="14">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(comsvcs\.dll.{0,40}MiniDump|rundll32.{0,40}comsvcs|procdump(64)?(\.exe)?.{0,40}lsass|-ma\s+lsass|lsass\.dmp|dumpert|nanodump)</field>
    <description>[T1003.001] LSASS memory dump command detected: $(win.eventdata.commandLine)</description>
    <mitre>
      <id>T1003</id>
      <id>T1003.001</id>
    </mitre>
    <group>credential_access,</group>
  </rule>

  <!-- =========================================================
    REGOLA 100034 — Keyword note di Mimikatz nella command line
    Trigger: Sysmon EID 1, comandi/moduli tipici di Mimikatz
    Tecnica: T1003 / T1003.001 (Credential Dumping)
    Severity: 15 (Critical) — quasi sempre malevolo
  ========================================================= -->
  <rule id="100034" level="15">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(sekurlsa::|logonpasswords|lsadump::|kerberos::|crypto::|privilege::debug|mimikatz|invoke-mimikatz|DumpCreds)</field>
    <description>[T1003] Mimikatz command-line keywords detected: $(win.eventdata.commandLine)</description>
    <mitre>
      <id>T1003</id>
      <id>T1003.001</id>
    </mitre>
    <group>credential_access,</group>
  </rule>

</group>
